# packaging/systemd/openstock-daily-pipeline.service
# openstock-daily-pipeline.service
# Runs the daily data pipeline via openstock-run-pipeline wrapper.
# Controlled by openstock-daily-pipeline.timer.
#
# Manual run:  sudo systemctl start openstock-daily-pipeline
# Logs:        journalctl -u openstock-daily-pipeline -f

[Unit]
Description=OpenStock daily data pipeline
Requires=openstock-data-platform.service
After=openstock-data-platform.service

[Service]
Type=oneshot
WorkingDirectory=/opt/openstock
EnvironmentFile=-/etc/openstock/openstock.env
ExecStart=/usr/local/bin/openstock-run-pipeline
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openstock-pipeline

[Install]
WantedBy=multi-user.target
