# packaging/systemd/openstock-provisioner.service
[Unit]
Description=OpenStock durable provisioning worker
Wants=network-online.target openstock-data-platform.service
After=network-online.target openstock-data-platform.service

[Service]
Type=simple
WorkingDirectory=/var/lib/openstock
EnvironmentFile=-/etc/vnalpha/vnalpha.env
Environment=VNALPHA_LOG_ROOT=/var/log/openstock
Environment=VNALPHA_PROVISIONING_QUEUE_PATH=/var/lib/openstock/queue/provisioning.sqlite3
Group=openstock
UMask=0007
# Retain the queue lock for the complete worker lifetime. Backup and restore
# acquire the same lock before making a paired warehouse/queue snapshot.
ExecStart=/usr/bin/flock -n /run/openstock-provisioner.lock /usr/bin/vnalpha provision worker --worker-id openstock-provisioner
Restart=on-failure
RestartSec=10s
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/openstock /var/log/openstock /run
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openstock-provisioner

[Install]
WantedBy=multi-user.target
