# packaging/deb/usr/lib/systemd/system/openstock-daily-pipeline.timer
[Unit]
Description=OpenStock weekday daily research-data maintenance timer

[Timer]
OnCalendar=Mon..Fri *-*-* 17:30:00 Asia/Ho_Chi_Minh
Persistent=true
RandomizedDelaySec=5min
AccuracySec=1min
Unit=openstock-daily-pipeline.service

[Install]
WantedBy=timers.target
